# base-repo/unit.d/etc-profile-d/overlay/etc/profile.d/10-build-paths.sh
# build path
export $(/usr/local/bin/path-helper)
