# base-repo/unit.d/etc-profile-d/overlay/etc/profile.d/50-kube-search-domains-fix.sh
# add kube search domains if set
if [[ $KUBE_SEARCH_DOMAINS_ADD == true ]]; then
    kube-search-domains add
fi

# remove kube search domains if set
if [[ $KUBE_SEARCH_DOMAINS_REMOVE == true ]]; then
    kube-search-domains remove
fi
