# base-repo/unit.d/etc-profile-d/overlay/etc/profile.d/50-kube-hostname-fix.sh
# activate hostname fix
if [[ $KUBE_HOSTNAME_FIX = true ]]; then
    export $(kube-hostname-fix enable)
fi
