# base-repo/unit.d/etc-profile-d/overlay/etc/profile.d/80-auto-detect-erlang-cookie.sh
if [ -s "${ERLANG_COOKIE_SRC:=/volumes/secrets/.erlang.cookie}" ] ||
    [ -n "$ERLANG_COOKIE" ]; then
        erlang-cookie-tool "auto-detect"
fi

