# base-repo/unit.d/etc-profile-d/overlay/etc/profile.d/80-fix-attrs.sh
# apply file ownership and permission fixes
if [[ ! $- =~ i ]]; then
    fixattrs
fi
